# package/debian/DEBIAN/prerm
#!/bin/bash
systemctl stop sensibo.service
systemctl disable sensibo.service
